# gcp-vm-image/packer/ansible/roles/base/files/etc/cron.d/rotate-opelog
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root

00 5 * * * root find /var/log/operation -mtime +3 -name "*.log" | xargs rm -f